# scripts/gh-release-proxy.service
[Unit]
Description=StrawWU GitHub Release CORS proxy (browser ISO join)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/node /mnt/data/code/project/StrawCoding/strawwu-public/scripts/gh-release-proxy.mjs
Environment=GH_RELEASE_PROXY_PORT=3911
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
